# roles/control_node/templates/route53-update.service.j2
[Unit]
Description = Run Route53 update jobs after Tower is up
After = ansible-tower.target

[Service]
ExecStart = /usr/local/bin/route53-update.py

[Install]
WantedBy = multi-user.target